# systemd/android-emulator-watchdog.service
[Unit]
Description=Android Emulator Farm Watchdog
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/android-emulator-farm
Environment="ANDROID_HOME=/opt/android-sdk"
Environment="ANDROID_SDK_ROOT=/opt/android-sdk"
Environment="PATH=/opt/android-sdk/emulator:/opt/android-sdk/platform-tools:/opt/android-sdk/cmdline-tools/latest/bin:/usr/bin:/bin"
ExecStart=python3 /opt/android-emulator-farm/scripts/emulator_manager.py watch
ExecStartPre=python3 /opt/android-emulator-farm/scripts/emulator_manager.py start
ExecStop=python3 /opt/android-emulator-farm/scripts/emulator_manager.py stop
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
