# systemd/appium-server.service
[Unit]
Description=Appium Server for Android Test Automation Farm
[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/appium

Environment=ANDROID_HOME=/root/Android/Sdk
Environment=ANDROID_SDK_ROOT=/root/Android/Sdk
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64
Environment=PATH=/root/.nvm/versions/node/v24.18.0/bin:/usr/local/bin:/usr/bin:/bin:/root/Android/Sdk/emulator:/root/Android/Sdk/platform-tools:/root/Android/Sdk/cmdline-tools/latest/bin

ExecStart=/root/.nvm/versions/node/v24.18.0/bin/appium server --config /opt/android-emulator-farm/config/appium-config.json

Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
