# systemd/appium-server.service
[Unit]
Description=Appium Server for Android Test Automation Farm
[Service]
Type=simple
User=appium
Group=appium
WorkingDirectory=/opt/appium

Environment=ANDROID_HOME=/opt/android-sdk
Environment=ANDROID_SDK_ROOT=/opt/android-sdk
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/opt/android-sdk/emulator:/opt/android-sdk/platform-tools:/opt/android-sdk/cmdline-tools/latest/bin

ExecStart=/usr/local/bin/appium server --config /opt/android-emulator-farm/config/appium-config.json

Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
